# docker/prod_vpp_agent/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug
nodaemon=true

[program:vpp]
command=/usr/bin/vpp -c /etc/vpp/vpp.conf
autorestart=false
redirect_stderr=true
priority=1

[program:agent]
command=/bin/vpp-agent --etcdv3-config=/opt/vpp-agent/dev/etcd.conf --kafka-config=/opt/vpp-agent/dev/kafka.conf --delay=80
autorestart=true
redirect_stderr=true
priority=2

; This event listener waits for event of vpp exitting. Once received, it kills supervisord process and this makes
; subsequently the exit of docker container.
; [eventlistener:vpp_not_running]
; command=/usr/bin/supervisord_kill.py
; events=PROCESS_STATE_EXITED
